# system_files/desktop/shared/etc/profile.d/kyawthuite-fastfetch
#!/usr/bin/env bash

# map fastfetch/neofetch to fastfetch with custom config
alias fastfetch='/usr/bin/fastfetch -c /usr/share/ublue-os/kyawthuite/fastfetch.jsonc'
alias neofetch='/usr/bin/fastfetch -c /usr/share/ublue-os/kyawthuite/fastfetch.jsonc'
