# system_files/desktop/shared/etc/profile.d/kyawthuite-fastfetch
alias fastfetch='/usr/bin/fastfetch -c /usr/share/ublue-os/kyawthuite/fastfetch.jsonc'
