# system-files/base/etc/systemd/system/dms-greeter-sync.service
[Unit]
Description=Copy DMS configs & wallpaper to global greeter cache

[Service]
Type=oneshot
ExecStart=/usr/libexec/dms-greeter-sync.sh
