# system_files/base/etc/profile.d/_sys.sh
#!/usr/bin/env bash

alias sys-update="_sys check"
alias sys-upgrade="_sys upgrade"
alias sys-upgrade-reboot="_sys reboot"