# system_files/base/etc/profile.d/fastfetch.sh
#!/usr/bin/env bash

alias neofetch="/usr/bin/fastfetch --logo /usr/share/icons/logo.txt --config /usr/share/ublue-os/kyawthuite/fastfetch.jsonc"

alias fastfetch="/usr/bin/fastfetch --logo /usr/share/icons/logo.txt --config /usr/share/ublue-os/kyawthuite/fastfetch.jsonc"