# system_files/base/etc/profile.d/sys.sh
#!/usr/bin/env bash

alias sys-update="sys check"
alias sys-upgrade="sys upgrade"
alias sys-upgrade-reboot="sys reboot"