# system-files/base/etc/profile.d/fastfetch.sh
#!/usr/bin/env bash

alias neofetch='/usr/bin/fastfetch --logo /usr/share/icons/logo.txt'
alias fastfetch='/usr/bin/fastfetch --logo /usr/share/icons/logo.txt'