# system_files/shared/etc/profile.d/askpass.sh
#!/usr/bin/env bash

# set `ksshaskpass` as default command for requesting SUDO pass
SUDO_ASKPASS='/usr/bin/ksshaskpass'

export SUDO_ASKPASS