# system_files/shared/etc/profile.d/open
#!/usr/bin/env bash

# use `open` command to open files from terminal
alias open="xdg-open &>/dev/null"